# templates/acme_renew.service.erb
[Unit]
Description=Run ACME renewals
After=network.target

[Service]
Type=simple
ExecStart=<%= @path %>/renew_all

