# files/acme_renew.service
[Unit]
Description=Run ACME renewals
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/certs/renew_all

